# packaging/sndctl.service
[Unit]
Description=Sound Control
Documentation=https://github.com/danmcpherson/Sndctl
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/sndctl
ExecStart=/opt/sndctl/.venv/bin/python -m uvicorn sndctl.main:app --host 0.0.0.0 --port 80
Restart=always
RestartSec=10
Environment=SNDCTL_DATA_DIRECTORY=/opt/sndctl/data
Environment=SNDCTL_WWWROOT_PATH=/opt/sndctl/wwwroot
Environment=SNDCTL_SOCO_CLI_USE_LOCAL_CACHE=false

# Allow binding to port 80
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
