# packaging/sndctl-upgrade.timer
[Unit]
Description=Sound Control Daily Upgrade Check Timer
Documentation=https://github.com/danmcpherson/Sndctl

[Timer]
# Run at 4am local time daily
# The exact time is randomized by up to 30 minutes to avoid thundering herd
OnCalendar=*-*-* 04:00:00
RandomizedDelaySec=1800
Persistent=true

# Start timer 5 minutes after boot if we missed the scheduled time
OnBootSec=5min

[Install]
WantedBy=timers.target
